# crontab
* * * * * /usr/local/bin/python /app/sync_script.py >> /var/log/cron.log 2>&1

